# userland/etc/mycelium/network.target
[Unit]
Description=network is configured
Requires=network.service

[Install]
WantedBy=multi-user.target
